# systemd/ryzen-power.path
[Unit]
Description=Trigger ryzen-power on platform_profile change

[Path]
PathChanged=/sys/firmware/acpi/platform_profile
Unit=ryzen-power.service

[Install]
WantedBy=multi-user.target
