# systemd/user/conditional-suspend.service
[Unit]
Description=Conditional suspend check (idle + SSH)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/conditional-suspend.sh

[Install]
WantedBy=default.target
